# keepalive-hub/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
(startup-notify.sh 2>&1 | logger) &
exit 0
