# elch/files/etc/uci-defaults/00-install-pkgs.sh
#!/bin/sh
opkg install --nodeps /pkgs/*
