# elch/files/etc/uci-defaults/00-install-pkgs.sh
#!/bin/sh

( sleep 10 && opkg update && opkg install /pkgs/*) &
