# elch/files/etc/uci-defaults/00-set-hostname.sh
#!/bin/sh
uci set system.@system[0].hostname=elch_$(hostid | md5sum | cut -c1-3 )
uci commit
