# elch/files/etc/uci-defaults/10-disable-services.sh
#!/bin/sh

## config dnsmasq
##   option rebind_protection '0'
# or:
/etc/init.d/dnsmasq stop
/etc/init.d/dnsmasq disable
