# elch/files/etc/uci-defaults/80-enable-services.sh
#!/bin/sh

/etc/init.d/uhttpd enable
/etc/init.d/collectd enable
