# keepalive-hub/backup/files/etc/uci-defaults/98-keepalive-backup.sh
#!/bin/sh
# TODO: configure to be backup (prio to 100), maybe even change the color - use green blinky
# for master and red blinky for backup
