# keepalive-hub/master/files/etc/uci-defaults/98-keepalive-master.sh
#!/bin/sh
# TODO: configure to be master, maybe even change the color - use green blinky
# for master and red blinky for backup
